# duckdns/duckdns.timer
[Unit]
Description=Run DuckDNS update script every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Unit=duckdns.service

[Install]
WantedBy=timers.target
